# modules/core/files/etc/profile.d/puppet.sh
#!/bin/bash
# puppet related aliases/functions

alias par='puppet apply --verbose --modulepath /etc/puppet/modules/ /etc/puppet/manifests/site.pp'